# pyinfra/red/files/snapraid-scrub.timer
[Unit]
Description=SnapRAID monthly scrub (10%)

[Timer]
OnCalendar=*-*-01 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
